# presentation/03-run-systemd/shbr-echo/shbr-echo.service
[Unit]
Description=shbr-echo
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/shbr-echo

[Install]
WantedBy=multi-user.target